# alertmanager-telegram-bridge.service
[Unit]
Description=Alertmanager Telegram Bridge
Documentation=https://github.com/bibigon14/alertmanager-telegram-bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bibigon88
Group=bibigon88

WorkingDirectory=/opt/alertmanager-telegram-bridge
ExecStart=/opt/alertmanager-telegram-bridge/venv/bin/python bridge.py

# Config path — override if needed
Environment="BRIDGE_CONFIG=/opt/alertmanager-telegram-bridge/config.yaml"

Restart=on-failure
RestartSec=5s

# Harden
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/alertmanager-telegram-bridge
PrivateTmp=true

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alertmanager-telegram-bridge

[Install]
WantedBy=multi-user.target
